# readyup.service
[Unit]
Description=ReadyUp Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=youruser   # ← change
Group=youruser  # ← change
WorkingDirectory=/opt/ready-up
ExecStart=/opt/ready-up/venv/bin/python /opt/ready-up/src/main.py

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

# ── security hardening ───────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
